# wordpress-nginx-nocache-xdebug/supervisor_fpm.conf
[program:fpm]
command=php5-fpm --nodaemonize
